# 00_System/launchd/com.leon.ip.process-ideas.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.leon.ip.process-ideas</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/anaconda3/bin/python3</string>
    <string>/Users/yige/Desktop/IP/00_System/scripts/ipctl.py</string>
    <string>process-ideas</string>
    <string>--execute</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yige/Desktop/IP</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/yige/.nvm/versions/node/v24.14.0/bin:/opt/anaconda3/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>10</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>19</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/yige/Desktop/IP/00_System/logs/launchd-process-ideas.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/yige/Desktop/IP/00_System/logs/launchd-process-ideas.err.log</string>
</dict>
</plist>
